# exampleCrontab.crontab
00/5 09-15 * * 1-5 /Users/amit/Projects/sayTime/say.sh
